# etc/crontab
0/30 * * * * /var/smoothwall/mods/redipnotify/red-change-notify.pl > /dev/null
